# .docker/services.d/php-fpm/run
#!/bin/execlineb -P
php-fpm -F -y /etc/php7/php-fpm.conf
